# conf/mysite.ini
[program:mysite]
command=python3 /root/Mysite/manage.py runserver 0.0.0.0:8080
stdoud_logfile=/opt/logs/mysite/django.log
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
